# manyfold/rootfs/etc/services.d/manyfold/run
#!/usr/bin/env bash
# ==============================================================================
# Home Assistant Community Add-on: Manyfold
# Runs the Manyfold application
# ==============================================================================

set -e

echo "Starting Manyfold..."

# Load configuration from Home Assistant
CONFIG_PATH="/data/options.json"

# Generate secret key if not exists
if [ ! -f /data/secret_key_base ]; then
    echo "Generating SECRET_KEY_BASE..."
    openssl rand -hex 64 > /data/secret_key_base
fi

# Export environment variables
export SECRET_KEY_BASE=$(cat /data/secret_key_base)
export DATABASE_ADAPTER="sqlite3"
export DATABASE_NAME="/data/manyfold.db"
export PUID=$(jq -r '.puid // 1000' ${CONFIG_PATH})
export PGID=$(jq -r '.pgid // 1000' ${CONFIG_PATH})

# Optional multiuser mode
if [ "$(jq -r '.multiuser_enabled // false' ${CONFIG_PATH})" = "true" ]; then
    export MULTIUSER="enabled"
    echo "Multiuser mode enabled"
    
    if [ "$(jq -r '.registration_enabled // false' ${CONFIG_PATH})" = "true" ]; then
        export REGISTRATION="enabled"
        echo "Registration enabled"
    fi
fi

# Optional usage tracking
if [ "$(jq -r '.usage_tracking_enabled // false' ${CONFIG_PATH})" = "true" ]; then
    export USAGE_TRACKING="enabled"
    echo "Usage tracking enabled"
fi

echo "Manyfold starting on port 3214..."
echo "Access at: http://homeassistant.local:3214"

# Execute the original Manyfold entrypoint
exec /init
